# utils/bw2_nap.service
[Unit]
Description="BW2 nuke and pave"

[Service]
Type=oneshot
ExecStart=/firmware/bw2_nap.sh
