# utils/fakehwclock.service
[Unit]
Description=Fake HW Clock save

[Service]
Type=simple
ExecStart=/opt/fake-hwclock.sh
