# utils/mfs.service
[Unit]
Description="Merged File System"

[Service]
Type=oneshot
RemainAfterExit=yes
StartLimitAction=reboot-force
ExecStart=/opt/mfs.sh
