# utils/prepwan.service
[Unit]
Description="Prep WAN"
After=mfs.service
Requires=mfs.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/firmware/prep_wan.sh
